# docker-supervisord.conf
[program:funding]
command=docker-compose up
directory=/code/fundingoss.com
redirect_stderr=true
autostart=true
autorestart=true
priority=10
